# work/config/includes.chroot_after_packages/etc/systemd/system/ndi-receiver.service
[Unit]
Description=NDI Receiver (SDL2)
After=network.target

[Service]
ExecStart=/path/to/ndi-receiver/start.sh
WorkingDirectory=/opt/ndi-receiver
Restart=always
RestartSec=3
Environment=XDG_RUNTIME_DIR=/run/user/%U

[Install]
WantedBy=default.target